# dev/docker/supervisord-sshd.conf
[program:sshd]
command=/start-sshd.sh
numprocs=1
autostart=true
autorestart=true